# system-files/btrfs-snabak.service
[Unit]
Description=Btrfs snapshot and backup daemon

[Service]
Type=simple
ExecStart=/usr/bin/btrfs-snabakd

[Install]
WantedBy=multi-user.target

